# gpio-agent/gpio-agent.service
[Unit]
Description=RoboFrame GPIO agent (PIR + 4x4 keypad + OLED)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Runs as root because evdev's UInput, gpiozero's MotionSensor on /dev/gpiomem
# and DigitalOutputDevice on /dev/gpiomem all need elevated access. Pi-only
# device; the security boundary is the Pi itself.
User=root
Environment=PYTHONUNBUFFERED=1
Environment=GPIO_AGENT_CONFIG=/opt/gpio-agent/config.json
WorkingDirectory=/opt/gpio-agent
ExecStart=/usr/bin/python3 /opt/gpio-agent/agent.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
